# custom-rpi-image/ext_dir/image/mbr/simple_dual/root-partition-resize-conf/expand-rootfs.service
[Unit]
Description=Expand root filesystem on first boot
After=local-fs.target
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/expand-rootfs.sh
RemainAfterExit=yes
User=root
Group=root

[Install]
WantedBy=multi-user.target